# annunciator/crontab
@reboot sleep 30 && /home/pi/OSIDS/annunciator/start.sh
